# scripts/systemd/podium-backend.service
[Unit]
Description=Podium relay + live agent daemon (backend, :18787)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user/src/other/podium
Environment=HOME=/home/user
Environment=PATH=/home/user/.local/bin:/home/user/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=PODIUM_PORT=18787
# Run @podium/* from TypeScript SOURCE (--conditions=@podium/source), exactly like
# Vite does (apps/web/vite.config.ts) — no build step, no dist, no stale-dist trap.
# Node runtime (via tsx) so the daemon's node-pty native addon loads (never Bun).
# No `watch`: a crash now EXITS so Restart=always actually recovers it (tsx watch
# would wedge a dead child while systemd still reported "active"). New main is picked
# up by podium-redeploy.path, which restarts this unit when .git/logs/HEAD moves.
ExecStart=/home/user/src/other/podium/node_modules/.bin/tsx --conditions=@podium/source scripts/host.ts
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
