# scripts/systemd/podium-health.timer
[Unit]
Description=Probe Podium backend health every 45s (last-resort wedge recovery)
# The 45s cadence is deliberately tight so a genuinely wedged /health is caught
# quickly — the timer is NOT what prevents false kills. The in-script guards
# (ActiveState check, 120s boot-age grace, two-probe confirmation) are what make
# a false restart structurally impossible; see scripts/podium-health-probe.sh
# and podium-health.service for the full rationale (POD-608).

[Timer]
OnBootSec=45s
OnUnitActiveSec=45s
AccuracySec=5s

[Install]
WantedBy=timers.target
