# scripts/systemd/podium-redeploy.service
[Unit]
Description=Podium redeploy — restart server + daemon + web to run the latest main (triggered by git HEAD change)
# Triggered by podium-redeploy.path on .git/logs/HEAD changes. The server + daemon run
# from source, so a restart picks up newly-merged code. The web Vite dev server must
# restart too: its long-lived module graph has gone stale across merges twice
# (stale @podium/protocol 2026-06-11; stale SocketHub/"hub.onAttention is not a
# function" 2026-06-12) — HMR does not reliably survive large git-driven rewrites.
After=podium-server.service podium-daemon.service

[Service]
Type=oneshot
# Close the redeploy/checkout race: a HEAD reflog write (our trigger) can land before
# git finishes the working-tree checkout, so wait for git to release .git/index.lock
# (the working tree is fully written then) before restarting. Bounded so a stuck lock
# can't wedge the redeploy.
ExecStartPre=/usr/bin/env bash /home/user/src/other/podium/scripts/redeploy-wait.sh /home/user/src/other/podium
# server first (the daemon reconnects to it with backoff); durable abduco/tmux masters
# survive the daemon restart in their own scopes, so live agents are not dropped.
ExecStart=/usr/bin/systemctl --user restart podium-server.service podium-daemon.service podium-web.service
